# deploy/borg.service
[Unit]
Description=Borg AI Agent Orchestrator
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/opt/borg
ExecStart=/opt/borg/borg-rs/target/release/borg-server
Restart=always
RestartSec=5
Environment=HOME=/root
Environment=PATH=/root/.cargo/bin:/root/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/bin
EnvironmentFile=/opt/borg/.env

[Install]
WantedBy=multi-user.target
